# deploy/systemd/blockrunnooor.service
[Unit]
Description=Blockrunnooor Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blockrun
Group=blockrun
WorkingDirectory=/opt/blockrunnooor

EnvironmentFile=/etc/blockrunnooor/blockrunnooor.env

ExecStart=/usr/bin/env bash -lc '/opt/blockrunnooor/bin/orchestrator'

Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/blockrunnooor
UMask=0077

[Install]
WantedBy=multi-user.target

